# backsched-check.timer
[Unit]
Description=Run backsched check every 24 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=24h

[Install]
WantedBy=default.target
