# backsched-check.service
[Unit]
Description=Check backup needed

[Service]
Type=oneshot
ExecStart=%h/bin/backsched need-backup --notify-send
