# backsched-check.service
[Unit]
Description=Check whether a backup needs to be performed.

[Service]
Type=oneshot
ExecStart=%h/bin/backsched check --notify
